# web2py/applications/Sensors/cron/crontab
#crontab
*/1 * * * * root python *applications/Sensors/modules/background_updater.py